# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --login-program '/usr/bin/chvt' --login-options 2 --noclear %I $TERM
Type=oneshot
TTYVTDisallocate=no
Restart=no